# deploy/roles/frontend/templates/steve-web.service.j2
# Managed by Ansible (frontend role).
# Next.js UI (withEve + useEveAgent). Caddy sends Eve protocol routes directly
# to the agent and sends all other paths here.
[Unit]
Description=steve-web — Next.js UI for the eve agent
After=network-online.target steve.service
Wants=network-online.target

[Service]
Type=simple
User={{ app_user }}
Group={{ app_user }}
WorkingDirectory={{ app_dir }}
Environment=HOME=/home/{{ app_user }}
Environment=NODE_ENV=production
Environment=NEXT_TELEMETRY_DISABLED=1
Environment=EVE_SELF_HOSTED=1
Environment=PORT={{ web_port }}
Environment=HOSTNAME=127.0.0.1
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/usr/bin/pnpm start --port {{ web_port }} --hostname 127.0.0.1
Restart=always
RestartSec=5
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target
