# deploy/roles/app/templates/steve.service.j2
# Managed by Ansible (app role).
# Long-running eve host. Uses `eve start` (the production host). As of eve 0.15.0
# `eve start` runs the custom Postgres world correctly (the old 0.13.x "Unhandled
# queue" regression is fixed). `eve build` runs in the app role before this unit
# starts, producing .output/ that `eve start` serves. Note: unlike `eve dev`,
# `eve start` does NOT reap per-run Docker sandbox containers on shutdown — run an
# external reaper if that matters for your deployment.
[Unit]
Description=steve — self-hosted eve agent host
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
User={{ app_user }}
Group={{ app_user }}
SupplementaryGroups=docker
WorkingDirectory={{ app_dir }}
EnvironmentFile={{ app_dir }}/.env
Environment=HOME=/home/{{ app_user }}
Environment=PORT={{ app_port }}
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
# Relax pnpm 11 supply-chain policies for the implicit deps-status check that
# `pnpm exec` runs (the PoC pins day-old beta versions on purpose).
Environment=PNPM_CONFIG_MINIMUM_RELEASE_AGE=0
Environment=PNPM_CONFIG_DANGEROUSLY_ALLOW_ALL_BUILDS=true
# Ensure the Postgres world container is up and accepting connections before
# the agent starts. On reboot the agent would otherwise race the Docker
# container and die with ECONNREFUSED to Postgres.
ExecStartPre=/usr/bin/docker compose -f {{ app_dir }}/docker-compose.yml up -d postgres
ExecStartPre=/bin/sh -c 'for i in $(seq 1 30); do docker exec steve-postgres pg_isready -U world -d world && exit 0; sleep 2; done; echo "postgres not ready" >&2; exit 1'
# Bind to localhost only; Caddy is the sole public entry point.
ExecStart=/usr/bin/pnpm exec eve start --host 127.0.0.1 --port {{ app_port }}
Restart=always
RestartSec=5
# Give the host time to pull the sandbox image on first boot.
TimeoutStartSec=300
KillMode=mixed

[Install]
WantedBy=multi-user.target
